# scripts/podman/autocrab.container.in
# AutoCrab gateway — Podman Quadlet (rootless)
# Installed by setup-podman.sh into autocrab's ~/.config/containers/systemd/
# {{AUTOCRAB_HOME}} is replaced at install time.

[Unit]
Description=AutoCrab gateway (rootless Podman)

[Container]
Image=autocrab:local
ContainerName=autocrab
UserNS=keep-id
# Keep container UID/GID aligned with the invoking user so mounted config is readable.
User=%U:%G
Volume={{AUTOCRAB_HOME}}/.autocrab:/home/node/.autocrab
EnvironmentFile={{AUTOCRAB_HOME}}/.autocrab/.env
Environment=HOME=/home/node
Environment=TERM=xterm-256color
PublishPort=18789:18789
PublishPort=18790:18790
Pull=never
Exec=node dist/index.js gateway --bind lan --port 18789

[Service]
TimeoutStartSec=300
Restart=on-failure

[Install]
WantedBy=default.target
